# conf/antenna.ini
[program:antenna-server]
directory = /antenna
command = python3 manage.py runserver 0.0.0.0:80
autostart = true
autorestart = true
redirect_stderr = true
stderr_logfile = /tmp/dnslog_stderr.log
stdout_logfile = /tmp/dnslog_stdout.log
stopsignal = KILL
stopasgroup = true
