# conf/antenna.ini
[program:antenna-server]
directory = /Antenna
command = python3 manage.py runserver 0.0.0.0:80
autostart = true
autorestart = true
redirect_stderr = true
stderr_logfile = /tmp/antenna_server_stderr.log
stdout_logfile = /tmp/antenna_server_stdout.log
stopsignal = KILL
stopasgroup = true

[program:antenna-dns]
directory = /Antenna
command = python3 modules/template/depend/listen/dnslog.py
autostart = true
autorestart = true
redirect_stderr = true
stderr_logfile = /tmp/antenna_dns_stderr.log
stdout_logfile = /tmp/antenna_dns_stdout.log
stopsignal = KILL
stopasgroup = true

[program:antenna-jndi]
directory = /Antenna
command = python3 modules/template/depend/listen/jndi.py
autostart = true
autorestart = true
redirect_stderr = true
stderr_logfile = /tmp/antenna_jndi_stderr.log
stdout_logfile = /tmp/antenna_jndi_stdout.log
stopsignal = KILL
stopasgroup = true

[program:antenna-ftp]
directory = /Antenna
command = python3 modules/template/depend/listen/ftplog.py
autostart = true
autorestart = true
redirect_stderr = true
stderr_logfile = /tmp/antenna_ftp_stderr.log
stdout_logfile = /tmp/antenna_ftp_stdout.log
stopsignal = KILL
stopasgroup = true

[program:antenna-https]
directory = /Antenna
command = python3 modules/template/depend/listen/httpslog.py
autostart = true
autorestart = true
redirect_stderr = true
stderr_logfile = /tmp/antenna_https_stderr.log
stdout_logfile = /tmp/antenna_https_stdout.log
stopsignal = KILL
stopasgroup = true

